# services/telegram-bot.service
[Unit]
Description=EasyClaw Telegram Bot Bridge
After=network-online.target claude-code.service
Wants=network-online.target

[Service]
Type=simple
User=%%USER%%
WorkingDirectory=%%EASYCLAW_DIR%%
EnvironmentFile=%%EASYCLAW_DIR%%/.env
ExecStart=/usr/bin/python3 %%EASYCLAW_DIR%%/telegram-bot.py
Restart=always
RestartSec=10s
StartLimitBurst=3
StartLimitIntervalSec=60


[Install]
WantedBy=multi-user.target
